# asu-cal.service
[Unit]
Description=ASU Sun Devil Athletics Calendar
After=network.target

[Service]
Type=simple
User=robert
WorkingDirectory=/home/robert/projects/asu-athletics-schedule
EnvironmentFile=-/home/robert/projects/secrets.env
ExecStart=/home/robert/.nvm/versions/node/v24.13.1/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
